# scripts/restflow.service
[Unit]
Description=RestFlow AI Agent Daemon
After=network.target

[Service]
Type=simple
User=restflow
Group=restflow
ExecStart=/usr/local/bin/restflow daemon start --foreground --http
ExecStop=/usr/local/bin/restflow daemon stop
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/restflow /var/log/restflow

[Install]
WantedBy=multi-user.target
